# operator/crontab
*/4 * * * * /opt/freeton-toolbox/operator/transactions-monitor.sh >> /var/log/freeton-toolbox/transactions-monitor.log 2>&1
