# operator/crontab
*/4 * * * * imapfilter >> /var/log/freeton-toolbox/imapfilter.log 2>&1
